# .docker/cron/periodic/daily/build-public
#!/bin/sh

curl -sS -H "X-Requested-With: XMLHttpRequest" -X POST http://$WEB_HOST/control-panel/build-private > /dev/stdout 2>/dev/stderr
